# debian/arsoft-xmpp-daemon.socket
[Socket]
ListenStream=/run/arsoft-xmpp-daemon/socket

[Install]
WantedBy=sockets.target
