# debian/arsoft-xmpp-daemon-setup.upstart
#!upstart
# arsoft-xmpp-daemon-setup

description     "AR Soft XMPP Daemon setup"
author          "Andreas Roth <aroth@arsoft-online.com>"

start on (local-filesystems and net-device-up IFACE!=lo and started udev-finish)

task
script
    [ ! -d /run/arsoft-xmpp-daemon ] && mkdir /run/arsoft-xmpp-daemon
    chown nobody.nogroup /run/arsoft-xmpp-daemon
    chmod 0755 /run/arsoft-xmpp-daemon
end script
#
# EOF
#
